# FPGA/Zynq_PS/rc.local
#!bin/sh -e
/home/NET2FPGA/bitstreamLoader.sh 
/home/NET2FPGA/constantsLoader 
exit 0
