# FPGA/Zynq_PS/rc.local
#!bin/sh -e
/home/NET2FPGA/setBitstream.sh 
/home/NET2FPGA/setConstants 
exit 0
