# Zynq_7010/Zynq_PS/rc.local
#!bin/sh -e
/home/Canvas/bitstreamLoader.sh 
/home/Canvas/constantsLoader 
exit 0
